# linux-clipboard-manager.service
[Unit]
Description=Linux Clipboard Manager
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/smart-clipboard
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
